# templates/default/supervisord.redis-queue-monitor.conf.erb
[program:redis_queue_monitor]
environment=PYTHONPATH=/var/rtb-devops/classes
command = /usr/bin/python /var/rtb-devops/remote-services/redis_queue_monitor.py

process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=KILL
stopasgroup = true
killasgroup = true


stdout_logfile=/tmp/supervisor_redis_queue_monitor.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10

stderr_logfile = /tmp/supervisor_redis_queue_monitor.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

[group:redis_queue_monitor_server]
programs=redis_queue_monitor






