# templates/default/supervisord.health-check.conf.erb
<% @role_list_hc.each do |role| %>



[program:<%="#{role}"%>]
environment=PYTHONPATH=/var/rtb-devops/classes
command = /usr/bin/python /var/rtb-devops/remote-services/datacenter_worker_healthcheck_solo.py --server_type <%="#{role}"%>


process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=KILL
stopasgroup = true
killasgroup = true


stdout_logfile=/tmp/supervisor_<%="#{role}"%>.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=1

stderr_logfile = /tmp/supervisor_<%="#{role}"%>.err
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=1


[group:<%="#{role}"%>_server]
programs=<%="#{role}"%>


<% end %>




