# templates/default/supervisord.ufw_firewall.conf.erb
[program:firewall]
environment=PYTHONPATH=/var/rtb-devops/classes
command = /usr/bin/python /var/rtb-devops/remote-services/ufw_firewall.py
autostart=true
autorestart=true

stopsignal=KILL
stopasgroup = true
killasgroup = true

stdout_logfile=/tmp/firewall.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10

stderr_logfile = /tmp/firewall.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

[group:firewall_server]
programs=firewall
