# templates/default/supervisord.coordinator.conf.erb
[program:coordinator]

command = /usr/bin/java -Xmx256m -Duser.timezone=UTC -Dfile.encoding=UTF-8 -classpath /var/druid/lib/*:/var/druid/config/coordinator io.druid.cli.Main server coordinator



process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=KILL
stopasgroup = true
killasgroup = true


stdout_logfile=/tmp/coordinator.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10

stderr_logfile = /tmp/coordinator.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10


[group:coordinator_server]
programs=coordinator




