# templates/default/supervisord.io-harvester.conf.erb
[program:harvester]
command = /usr/bin/node /usr/bin/log.io-harvester

process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=KILL
stopasgroup = true
killasgroup = true


stdout_logfile=/tmp/supervisor_harvester.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10

stderr_logfile = /tmp/supervisor_harvester.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

[group:harvester_server]
programs=harvester






